# scripts/atlas-firstboot.service
[Unit]
Description=ODS Player OS Atlas — Automated First Boot Setup
After=network-online.target
Wants=network-online.target
# Use ODS-owned gate file instead of Armbian's .not_logged_in_yet
# (Armbian's first-login deletes the old file before we can run — race condition)
ConditionPathExists=/var/lib/ods/atlas_firstboot_pending

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/atlas_firstboot.sh
# Output to tty1 (console) + journal so boot steps are visible on screen
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
RemainAfterExit=yes
Environment=HOME=/root
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
